# tinyfecvpn-easyconf/net/etc/hotplug.d/iface/30-tinyfec
ip route add $(uci get tinyvpnconf.Setup.dstAddr)/32 via $(ip route | awk 'NR==2{print $3}')